# root/etc/profile.d/o.sh
export PATH="$PATH:$HOME/.local/share/flatpak/exports/bin"

export EDITOR="/bin/nvim"
